# scripts/volumebegone.service
[Unit]
Description=Volume Be Gone - Bluetooth Speaker Control
After=bluetooth.target sound.target network.target
Wants=bluetooth.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/volumebegone
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /home/pi/volumebegone/volumeBeGone.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
